# .config/systemd/user/luncher.service
[Unit]
Description=luncher: custom luncher

[Service]
Type=simple
ExecStart=/home/cndul/bin/luncher

[Install]
WantedBy=default.target
