# .config/systemd/user/my-timer.service
[Unit]
Description=my timer

[Service]
Type=simple
ExecStart=/home/cndul/bin/take-break.sh
